# packages/pfsense-package/pkg-plist
etc/rc.d/monitor_pfsense_agent
libexec/monitor-pfsense-agent/monitor-pfsense-agent-loop.sh
libexec/monitor-pfsense-agent/monitor-pfsense-agent.sh
pkg/systemup_monitor.inc
pkg/systemup_monitor.xml
share/pfSense-pkg-systemup-monitor/info.xml
share/pfSense-pkg-systemup-monitor/systemup_monitor_cli.php
www/config_systemup_monitor.php
www/status_systemup_monitor.php
